# aws-solution/ansible/roles/node_exporter/templates/node_exporter.service.j2

[Unit]
Description=node_exporter - Exporter for machine metrics.
After=network.target

[Service]
ExecStart={{node_exporter_path}}/node_exporter -web.listen-address :32001
